# deploy/systemd/otterbot.service
[Unit]
Description=Otterbot Discord Bot
After=network.target

[Service]
Type=simple
User=botuser
WorkingDirectory=/home/botuser/projects/otterbot
Environment=NODE_ENV=production
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
MemoryMax=750M
MemorySwapMax=1500M
StandardOutput=journal
StandardError=journal
SyslogIdentifier=otterbot

[Install]
WantedBy=multi-user.target

# NOTE: If /usr/bin/node is not correct for your system, run:
#   which node
# and update ExecStart accordingly.
#
# To install/update this service on the VPS:
#   sudo cp deploy/systemd/otterbot.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable otterbot
#   sudo systemctl restart otterbot
