# deploy/systemd/otterbot-restart.service
[Unit]
Description=Weekly Otterbot restart

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart otterbot.service
